# Ubuntu_deb/weavedconnectd-1.3-06_x86/weavedconnectd-1.3-06_x86/etc/init/weaved.conf
description "Weaved Remot3.it"
author      "mike@weaved"

start on (local-filesystems and net-device-up)
stop on runlevel [!2345]

pre-start script
    logger "[weavedweaved] upstart "
    sleep 70
    logger "[weavedweaved] upstart 70"
end script

script
    logger "[weavedweaved] script"
    exec /etc/init.d/weaved restart
    exec /etc/init.d/weaved.schannel restart
end script

